# autocorrect/rootfs/etc/ha-autocorrect.d/s6-overlay/s6-rc.d/autocorrect/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start service
# ==============================================================================
extra_args=()
if bashio::config.true 'debug_logging'; then
    extra_args+=('--debug')
fi

if bashio::config.true 'sound_enabled'; then
    extra_args+=('--snd-command' 'aplay -r 16000 -c 1 -f S16_LE -t raw')
fi

exec /usr/bin/entrypoint.sh